# OpenRC-Init-Script/cups-browsed
#!/sbin/openrc-run

name="CUPS Browsed"
description="Make remote CUPS printers available locally"
command="/usr/bin/cups-browsed"
supervisor="supervise-daemon"
pidfile="/run/cups-browsed.pid"

depend() {
    need dbus
    after cupsd avahi-daemon net
}

start_pre() {
    checkpath -d -m 0755 -o root:lp /run/cups
}


